# docker/root/etc/services.d/fixmylib/run
#!/usr/bin/with-contenv bash

echo "Starting fixmylib..."

cd /app || exit 1

exec s6-setuidgid fixmylib ./fixmylib
